# etc/supervisor/conf.d/rabbitmq.conf
[program:rabbitmq]
command=/usr/bin/rabbitmq-wrapper.sh
